# configuration/etc/network/if-up.d/karoshi-make-sure-up
#!/bin/bash

#Wait until the interface is *actually* up
timeout 10 while true; do
	if [[ $(< /sys/class/net/$IFACE/operstate) == "up" ]]; then
		true
		exit
	fi
	sleep 0.1
done
