# clientsetup/buildclient/config/configuration/etc/xdg/autostart/karoshi-ibus.desktop
[Desktop Entry]
Type=Application
Name=Karoshi iBus launcher
OnlyShowIn=XFCE;
Exec=ibus-daemon
