# clientsetup/buildclient/config/configuration/etc/xdg/autostart/karoshi-pre-session.desktop
[Desktop Entry]
Type=Application
Name=Karoshi Pre session
OnlyShowIn=XFCE;
Exec=/opt/karoshi/clientsetup/scripts/pre-session
