# linuxclientsetup/utilities/karoshi-remount.conf
# Karoshi Remount
#
# This task remounts certain directories that can become world writable to
# apply nosuid and nodev as necessary

description "Karoshi Remount"
author      "Robin McCorkell <rmccorkell@karoshi.org.uk>"

start on (filesystem)

console log
task

script
	mount -o remount,nosuid /home
	mount -o remount,nosuid /var
	mount -o remount,nosuid,nodev /tmp
end script

